# docker/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisor.log      ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB                           ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                              ; (num of main logfile rotation backups;default 10)
loglevel=info                                   ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid                ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                                  ; (start in foreground if true;default false)
minfds=1024                                     ; (min. avail startup file descriptors;default 1024)
minprocs=200                                    ; (min. avail process descriptors;default 200)

[program:sshd]
command=/usr/sbin/sshd -D

[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true

#~ [program:gunicorn]
#~ command=/usr/sbin/gunicorn-debian --conf-dir=/etc/gunicorn.d/ -pid-dir=/var/run --log-dir=/var/log/gunicorn
#~ redirect_stderr=true

[include]
files = /etc/supervisor/janitoo.conf.d/*.conf
