# docker/supervisord.conf
[supervisord]
nodaemon=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf

